# deployment/systemd/crypto13-live-research.service.example
[Unit]
Description=Crypto13Research live paper research example
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/crypto13_research
Environment=API_MODE=paper
Environment=ALLOW_REAL_ORDERS=false
Environment=ALLOW_TESTNET_ORDERS=false
ExecStart=/usr/bin/python -m src.main live-research --tf 15m --interval-sec 60 --run-forever
Restart=on-failure

[Install]
WantedBy=multi-user.target
